# smart-eldercare.service
[Unit]
Description=智慧养老平台后端服务
After=network.target mysql.service

[Service]
Type=simple
User=root
WorkingDirectory=__PROJECT_DIR__
ExecStart=/bin/bash -c 'source __PROJECT_DIR__/.env && exec /usr/bin/java -jar __PROJECT_DIR__/target/smart-elder-care-0.0.1-SNAPSHOT.jar'
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
